# supervisor/benchmark.conf
[program:benchmark]
directory=/<path>/arewefastyet
command=/<path>/arewefastyet/benchmark/bin/gunicorn server:app -b localhost:8000
autostart=true
autorestart=true